# service/abcd.service.in
[Unit]
Description=ABCD data acquisition
After=network.target

[Service]
Type=oneshot
User={{ username }}
ExecStart=/bin/bash {{ startup_script }}
ExecStop=/usr/bin/tmux kill-session -t ABCD
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

